# pkgsrc.2022Q2/linux/audit/PLIST
@comment $NetBSD$
bin/aulast
bin/aulastlog
bin/ausyscall
bin/auvirt
etc/audit/audisp-remote.conf
etc/audit/audit-stop.rules
etc/audit/auditd.conf
etc/audit/plugins.d/af_unix.conf
etc/audit/plugins.d/au-remote.conf
etc/audit/plugins.d/syslog.conf
etc/libaudit.conf
etc/rc.d/init.d/auditd
etc/sysconfig/auditd
include/auparse-defs.h
include/auparse.h
include/libaudit.h
lib/libaudit.la
lib/libauparse.la
lib/pkgconfig/audit.pc
lib/pkgconfig/auparse.pc
man/man3/audit_add_rule_data.3
man/man3/audit_add_watch.3
man/man3/audit_close.3
man/man3/audit_delete_rule_data.3
man/man3/audit_detect_machine.3
man/man3/audit_encode_nv_string.3
man/man3/audit_encode_value.3
man/man3/audit_flag_to_name.3
man/man3/audit_fstype_to_name.3
man/man3/audit_get_reply.3
man/man3/audit_get_session.3
man/man3/audit_getloginuid.3
man/man3/audit_is_enabled.3
man/man3/audit_log_acct_message.3
man/man3/audit_log_semanage_message.3
man/man3/audit_log_user_avc_message.3
man/man3/audit_log_user_comm_message.3
man/man3/audit_log_user_command.3
man/man3/audit_log_user_message.3
man/man3/audit_name_to_action.3
man/man3/audit_name_to_errno.3
man/man3/audit_name_to_flag.3
man/man3/audit_name_to_fstype.3
man/man3/audit_name_to_syscall.3
man/man3/audit_open.3
man/man3/audit_request_rules_list_data.3
man/man3/audit_request_signal_info.3
man/man3/audit_request_status.3
man/man3/audit_set_backlog_limit.3
man/man3/audit_set_backlog_wait_time.3
man/man3/audit_set_enabled.3
man/man3/audit_set_failure.3
man/man3/audit_set_pid.3
man/man3/audit_set_rate_limit.3
man/man3/audit_setloginuid.3
man/man3/audit_syscall_to_name.3
man/man3/audit_update_watch_perms.3
man/man3/audit_value_needs_encoding.3
man/man3/auparse_add_callback.3
man/man3/auparse_destroy.3
man/man3/auparse_feed.3
man/man3/auparse_feed_age_events.3
man/man3/auparse_feed_has_data.3
man/man3/auparse_find_field.3
man/man3/auparse_find_field_next.3
man/man3/auparse_first_field.3
man/man3/auparse_first_record.3
man/man3/auparse_flush_feed.3
man/man3/auparse_get_field_int.3
man/man3/auparse_get_field_name.3
man/man3/auparse_get_field_num.3
man/man3/auparse_get_field_str.3
man/man3/auparse_get_field_type.3
man/man3/auparse_get_filename.3
man/man3/auparse_get_line_number.3
man/man3/auparse_get_milli.3
man/man3/auparse_get_node.3
man/man3/auparse_get_num_fields.3
man/man3/auparse_get_num_records.3
man/man3/auparse_get_record_num.3
man/man3/auparse_get_record_text.3
man/man3/auparse_get_serial.3
man/man3/auparse_get_time.3
man/man3/auparse_get_timestamp.3
man/man3/auparse_get_type.3
man/man3/auparse_get_type_name.3
man/man3/auparse_goto_field_num.3
man/man3/auparse_goto_record_num.3
man/man3/auparse_init.3
man/man3/auparse_interpret_field.3
man/man3/auparse_new_buffer.3
man/man3/auparse_next_event.3
man/man3/auparse_next_field.3
man/man3/auparse_next_record.3
man/man3/auparse_node_compare.3
man/man3/auparse_normalize.3
man/man3/auparse_normalize_functions.3
man/man3/auparse_reset.3
man/man3/auparse_set_eoe_timeout.3
man/man3/auparse_set_escape_mode.3
man/man3/auparse_timestamp_compare.3
man/man3/ausearch_add_expression.3
man/man3/ausearch_add_interpreted_item.3
man/man3/ausearch_add_item.3
man/man3/ausearch_add_regex.3
man/man3/ausearch_add_timestamp_item.3
man/man3/ausearch_add_timestamp_item_ex.3
man/man3/ausearch_clear.3
man/man3/ausearch_next_event.3
man/man3/ausearch_set_stop.3
man/man3/get_auditfail_action.3
man/man3/set_aumessage_mode.3
man/man5/audisp-remote.conf.5
man/man5/auditd-plugins.5
man/man5/auditd.conf.5
man/man5/ausearch-expression.5
man/man5/libaudit.conf.5
man/man5/zos-remote.conf.5
man/man7/audit.rules.7
man/man8/audisp-remote.8
man/man8/audisp-syslog.8
man/man8/audispd-zos-remote.8
man/man8/auditctl.8
man/man8/auditd.8
man/man8/augenrules.8
man/man8/aulast.8
man/man8/aulastlog.8
man/man8/aureport.8
man/man8/ausearch.8
man/man8/ausyscall.8
man/man8/autrace.8
man/man8/auvirt.8
native/bin/aulast
native/bin/aulastlog
native/bin/ausyscall
native/bin/auvirt
native/etc/audit/audisp-remote.conf
native/etc/audit/audit-stop.rules
native/etc/audit/auditd.conf
native/etc/audit/plugins.d/af_unix.conf
native/etc/audit/plugins.d/au-remote.conf
native/etc/audit/plugins.d/syslog.conf
native/etc/libaudit.conf
native/etc/rc.d/init.d/auditd
native/etc/sysconfig/auditd
native/sbin/audisp-remote
native/sbin/audisp-syslog
native/sbin/auditctl
native/sbin/auditd
native/sbin/augenrules
native/sbin/aureport
native/sbin/ausearch
native/sbin/autrace
native/usr/include/auparse-defs.h
native/usr/include/auparse.h
native/usr/include/libaudit.h
native/usr/lib/libaudit.a
native/usr/lib/libaudit.la
native/usr/lib/libaudit.so
native/usr/lib/libaudit.so.1
native/usr/lib/libaudit.so.1.0.0
native/usr/lib/libauparse.a
native/usr/lib/libauparse.la
native/usr/lib/libauparse.so
native/usr/lib/libauparse.so.0
native/usr/lib/libauparse.so.0.0.0
native/usr/lib/pkgconfig/audit.pc
native/usr/lib/pkgconfig/auparse.pc
native/usr/share/aclocal/audit.m4
native/usr/share/audit/sample-rules/10-base-config.rules
native/usr/share/audit/sample-rules/10-no-audit.rules
native/usr/share/audit/sample-rules/11-loginuid.rules
native/usr/share/audit/sample-rules/12-cont-fail.rules
native/usr/share/audit/sample-rules/12-ignore-error.rules
native/usr/share/audit/sample-rules/20-dont-audit.rules
native/usr/share/audit/sample-rules/21-no32bit.rules
native/usr/share/audit/sample-rules/22-ignore-chrony.rules
native/usr/share/audit/sample-rules/23-ignore-filesystems.rules
native/usr/share/audit/sample-rules/30-nispom.rules
native/usr/share/audit/sample-rules/30-ospp-v42-1-create-failed.rules
native/usr/share/audit/sample-rules/30-ospp-v42-1-create-success.rules
native/usr/share/audit/sample-rules/30-ospp-v42-2-modify-failed.rules
native/usr/share/audit/sample-rules/30-ospp-v42-2-modify-success.rules
native/usr/share/audit/sample-rules/30-ospp-v42-3-access-failed.rules
native/usr/share/audit/sample-rules/30-ospp-v42-3-access-success.rules
native/usr/share/audit/sample-rules/30-ospp-v42-4-delete-failed.rules
native/usr/share/audit/sample-rules/30-ospp-v42-4-delete-success.rules
native/usr/share/audit/sample-rules/30-ospp-v42-5-perm-change-failed.rules
native/usr/share/audit/sample-rules/30-ospp-v42-5-perm-change-success.rules
native/usr/share/audit/sample-rules/30-ospp-v42-6-owner-change-failed.rules
native/usr/share/audit/sample-rules/30-ospp-v42-6-owner-change-success.rules
native/usr/share/audit/sample-rules/30-ospp-v42.rules
native/usr/share/audit/sample-rules/30-pci-dss-v31.rules
native/usr/share/audit/sample-rules/30-stig.rules
native/usr/share/audit/sample-rules/31-privileged.rules
native/usr/share/audit/sample-rules/32-power-abuse.rules
native/usr/share/audit/sample-rules/40-local.rules
native/usr/share/audit/sample-rules/41-containers.rules
native/usr/share/audit/sample-rules/42-injection.rules
native/usr/share/audit/sample-rules/43-module-load.rules
native/usr/share/audit/sample-rules/44-installers.rules
native/usr/share/audit/sample-rules/70-einval.rules
native/usr/share/audit/sample-rules/71-networking.rules
native/usr/share/audit/sample-rules/99-finalize.rules
native/usr/share/audit/sample-rules/README-rules
native/usr/share/examples/rc.d/auditd
native/usr/share/man/man3/audit_add_rule_data.3
native/usr/share/man/man3/audit_add_watch.3
native/usr/share/man/man3/audit_close.3
native/usr/share/man/man3/audit_delete_rule_data.3
native/usr/share/man/man3/audit_detect_machine.3
native/usr/share/man/man3/audit_encode_nv_string.3
native/usr/share/man/man3/audit_encode_value.3
native/usr/share/man/man3/audit_flag_to_name.3
native/usr/share/man/man3/audit_fstype_to_name.3
native/usr/share/man/man3/audit_get_reply.3
native/usr/share/man/man3/audit_get_session.3
native/usr/share/man/man3/audit_getloginuid.3
native/usr/share/man/man3/audit_is_enabled.3
native/usr/share/man/man3/audit_log_acct_message.3
native/usr/share/man/man3/audit_log_semanage_message.3
native/usr/share/man/man3/audit_log_user_avc_message.3
native/usr/share/man/man3/audit_log_user_comm_message.3
native/usr/share/man/man3/audit_log_user_command.3
native/usr/share/man/man3/audit_log_user_message.3
native/usr/share/man/man3/audit_name_to_action.3
native/usr/share/man/man3/audit_name_to_errno.3
native/usr/share/man/man3/audit_name_to_flag.3
native/usr/share/man/man3/audit_name_to_fstype.3
native/usr/share/man/man3/audit_name_to_syscall.3
native/usr/share/man/man3/audit_open.3
native/usr/share/man/man3/audit_request_rules_list_data.3
native/usr/share/man/man3/audit_request_signal_info.3
native/usr/share/man/man3/audit_request_status.3
native/usr/share/man/man3/audit_set_backlog_limit.3
native/usr/share/man/man3/audit_set_backlog_wait_time.3
native/usr/share/man/man3/audit_set_enabled.3
native/usr/share/man/man3/audit_set_failure.3
native/usr/share/man/man3/audit_set_pid.3
native/usr/share/man/man3/audit_set_rate_limit.3
native/usr/share/man/man3/audit_setloginuid.3
native/usr/share/man/man3/audit_syscall_to_name.3
native/usr/share/man/man3/audit_update_watch_perms.3
native/usr/share/man/man3/audit_value_needs_encoding.3
native/usr/share/man/man3/auparse_add_callback.3
native/usr/share/man/man3/auparse_destroy.3
native/usr/share/man/man3/auparse_feed.3
native/usr/share/man/man3/auparse_feed_age_events.3
native/usr/share/man/man3/auparse_feed_has_data.3
native/usr/share/man/man3/auparse_find_field.3
native/usr/share/man/man3/auparse_find_field_next.3
native/usr/share/man/man3/auparse_first_field.3
native/usr/share/man/man3/auparse_first_record.3
native/usr/share/man/man3/auparse_flush_feed.3
native/usr/share/man/man3/auparse_get_field_int.3
native/usr/share/man/man3/auparse_get_field_name.3
native/usr/share/man/man3/auparse_get_field_num.3
native/usr/share/man/man3/auparse_get_field_str.3
native/usr/share/man/man3/auparse_get_field_type.3
native/usr/share/man/man3/auparse_get_filename.3
native/usr/share/man/man3/auparse_get_line_number.3
native/usr/share/man/man3/auparse_get_milli.3
native/usr/share/man/man3/auparse_get_node.3
native/usr/share/man/man3/auparse_get_num_fields.3
native/usr/share/man/man3/auparse_get_num_records.3
native/usr/share/man/man3/auparse_get_record_num.3
native/usr/share/man/man3/auparse_get_record_text.3
native/usr/share/man/man3/auparse_get_serial.3
native/usr/share/man/man3/auparse_get_time.3
native/usr/share/man/man3/auparse_get_timestamp.3
native/usr/share/man/man3/auparse_get_type.3
native/usr/share/man/man3/auparse_get_type_name.3
native/usr/share/man/man3/auparse_goto_field_num.3
native/usr/share/man/man3/auparse_goto_record_num.3
native/usr/share/man/man3/auparse_init.3
native/usr/share/man/man3/auparse_interpret_field.3
native/usr/share/man/man3/auparse_new_buffer.3
native/usr/share/man/man3/auparse_next_event.3
native/usr/share/man/man3/auparse_next_field.3
native/usr/share/man/man3/auparse_next_record.3
native/usr/share/man/man3/auparse_node_compare.3
native/usr/share/man/man3/auparse_normalize.3
native/usr/share/man/man3/auparse_normalize_functions.3
native/usr/share/man/man3/auparse_reset.3
native/usr/share/man/man3/auparse_set_eoe_timeout.3
native/usr/share/man/man3/auparse_set_escape_mode.3
native/usr/share/man/man3/auparse_timestamp_compare.3
native/usr/share/man/man3/ausearch_add_expression.3
native/usr/share/man/man3/ausearch_add_interpreted_item.3
native/usr/share/man/man3/ausearch_add_item.3
native/usr/share/man/man3/ausearch_add_regex.3
native/usr/share/man/man3/ausearch_add_timestamp_item.3
native/usr/share/man/man3/ausearch_add_timestamp_item_ex.3
native/usr/share/man/man3/ausearch_clear.3
native/usr/share/man/man3/ausearch_next_event.3
native/usr/share/man/man3/ausearch_set_stop.3
native/usr/share/man/man3/get_auditfail_action.3
native/usr/share/man/man3/set_aumessage_mode.3
native/usr/share/man/man5/audisp-remote.conf.5
native/usr/share/man/man5/auditd-plugins.5
native/usr/share/man/man5/auditd.conf.5
native/usr/share/man/man5/ausearch-expression.5
native/usr/share/man/man5/libaudit.conf.5
native/usr/share/man/man5/zos-remote.conf.5
native/usr/share/man/man7/audit.rules.7
native/usr/share/man/man8/audisp-remote.8
native/usr/share/man/man8/audisp-syslog.8
native/usr/share/man/man8/audispd-zos-remote.8
native/usr/share/man/man8/auditctl.8
native/usr/share/man/man8/auditd.8
native/usr/share/man/man8/augenrules.8
native/usr/share/man/man8/aulast.8
native/usr/share/man/man8/aulastlog.8
native/usr/share/man/man8/aureport.8
native/usr/share/man/man8/ausearch.8
native/usr/share/man/man8/ausyscall.8
native/usr/share/man/man8/autrace.8
native/usr/share/man/man8/auvirt.8
sbin/audisp-remote
sbin/audisp-syslog
sbin/auditctl
sbin/auditd
sbin/augenrules
sbin/aureport
sbin/ausearch
sbin/autrace
share/aclocal/audit.m4
share/audit/sample-rules/10-base-config.rules
share/audit/sample-rules/10-no-audit.rules
share/audit/sample-rules/11-loginuid.rules
share/audit/sample-rules/12-cont-fail.rules
share/audit/sample-rules/12-ignore-error.rules
share/audit/sample-rules/20-dont-audit.rules
share/audit/sample-rules/21-no32bit.rules
share/audit/sample-rules/22-ignore-chrony.rules
share/audit/sample-rules/23-ignore-filesystems.rules
share/audit/sample-rules/30-nispom.rules
share/audit/sample-rules/30-ospp-v42-1-create-failed.rules
share/audit/sample-rules/30-ospp-v42-1-create-success.rules
share/audit/sample-rules/30-ospp-v42-2-modify-failed.rules
share/audit/sample-rules/30-ospp-v42-2-modify-success.rules
share/audit/sample-rules/30-ospp-v42-3-access-failed.rules
share/audit/sample-rules/30-ospp-v42-3-access-success.rules
share/audit/sample-rules/30-ospp-v42-4-delete-failed.rules
share/audit/sample-rules/30-ospp-v42-4-delete-success.rules
share/audit/sample-rules/30-ospp-v42-5-perm-change-failed.rules
share/audit/sample-rules/30-ospp-v42-5-perm-change-success.rules
share/audit/sample-rules/30-ospp-v42-6-owner-change-failed.rules
share/audit/sample-rules/30-ospp-v42-6-owner-change-success.rules
share/audit/sample-rules/30-ospp-v42.rules
share/audit/sample-rules/30-pci-dss-v31.rules
share/audit/sample-rules/30-stig.rules
share/audit/sample-rules/31-privileged.rules
share/audit/sample-rules/32-power-abuse.rules
share/audit/sample-rules/40-local.rules
share/audit/sample-rules/41-containers.rules
share/audit/sample-rules/42-injection.rules
share/audit/sample-rules/43-module-load.rules
share/audit/sample-rules/44-installers.rules
share/audit/sample-rules/70-einval.rules
share/audit/sample-rules/71-networking.rules
share/audit/sample-rules/99-finalize.rules
share/audit/sample-rules/README-rules
