# pkglinux-pkgsrc.2022Q2/native/openrc/PLIST
@comment $NetBSD$
native/bin/rc-status
native/etc/conf.d/agetty
native/etc/conf.d/agetty.tty1
native/etc/conf.d/bootmisc
native/etc/conf.d/consolefont
native/etc/conf.d/devfs
native/etc/conf.d/dmesg
native/etc/conf.d/fsck
native/etc/conf.d/hostname
native/etc/conf.d/hwclock
native/etc/conf.d/keymaps
native/etc/conf.d/killprocs
native/etc/conf.d/localmount
native/etc/conf.d/modules
native/etc/conf.d/mtab
native/etc/conf.d/net-online
native/etc/conf.d/netmount
native/etc/conf.d/network
native/etc/conf.d/seedrng
native/etc/conf.d/staticroute
native/etc/conf.d/swap
native/etc/conf.d/swclock
native/etc/conf.d/udev
native/etc/conf.d/udev-settle
native/etc/conf.d/udev-trigger
native/etc/init.d/agetty
native/etc/init.d/agetty.tty1
native/etc/init.d/agetty.tty2
native/etc/init.d/agetty.tty3
native/etc/init.d/agetty.tty4
native/etc/init.d/agetty.tty5
native/etc/init.d/agetty.tty6
native/etc/init.d/binfmt
native/etc/init.d/bootmisc
native/etc/init.d/cgroups
native/etc/init.d/consolefont
native/etc/init.d/devfs
native/etc/init.d/dmesg
native/etc/init.d/fsck
native/etc/init.d/functions.sh
native/etc/init.d/hostname
native/etc/init.d/hwclock
native/etc/init.d/keymaps
native/etc/init.d/killprocs
native/etc/init.d/kmod-static-nodes
native/etc/init.d/local
native/etc/init.d/localmount
native/etc/init.d/loopback
native/etc/init.d/modules
native/etc/init.d/mount-ro
native/etc/init.d/mtab
native/etc/init.d/net-online
native/etc/init.d/netmount
native/etc/init.d/network
native/etc/init.d/nscd
native/etc/init.d/numlock
native/etc/init.d/osclock
native/etc/init.d/procfs
native/etc/init.d/root
native/etc/init.d/runsvdir
native/etc/init.d/s6-svscan
native/etc/init.d/save-keymaps
native/etc/init.d/save-termencoding
native/etc/init.d/savecache
native/etc/init.d/seedrng
native/etc/init.d/staticroute
native/etc/init.d/swap
native/etc/init.d/swclock
native/etc/init.d/sysctl
native/etc/init.d/sysfs
native/etc/init.d/termencoding
native/etc/init.d/udev
native/etc/init.d/udev-settle
native/etc/init.d/udev-trigger
native/etc/local.d/README
native/etc/local.d/rc.local.start
native/etc/local.d/rc.local.stop
native/etc/rc.conf
native/etc/runlevels/boot/binfmt
native/etc/runlevels/boot/bootmisc
native/etc/runlevels/boot/consolefont
native/etc/runlevels/boot/fsck
native/etc/runlevels/boot/hostname
native/etc/runlevels/boot/hwclock
native/etc/runlevels/boot/keymaps
native/etc/runlevels/boot/localmount
native/etc/runlevels/boot/loopback
native/etc/runlevels/boot/modules
native/etc/runlevels/boot/mtab
native/etc/runlevels/boot/network
native/etc/runlevels/boot/procfs
native/etc/runlevels/boot/root
native/etc/runlevels/boot/save-keymaps
native/etc/runlevels/boot/save-termencoding
native/etc/runlevels/boot/seedrng
native/etc/runlevels/boot/staticroute
native/etc/runlevels/boot/swap
native/etc/runlevels/boot/sysctl
native/etc/runlevels/boot/termencoding
native/etc/runlevels/default/agetty.tty1
native/etc/runlevels/default/agetty.tty2
native/etc/runlevels/default/agetty.tty3
native/etc/runlevels/default/agetty.tty4
native/etc/runlevels/default/agetty.tty5
native/etc/runlevels/default/agetty.tty6
native/etc/runlevels/default/local
native/etc/runlevels/default/netmount
native/etc/runlevels/default/nscd
native/etc/runlevels/nonetwork/local
native/etc/runlevels/shutdown/killprocs
native/etc/runlevels/shutdown/mount-ro
native/etc/runlevels/shutdown/savecache
native/etc/runlevels/sysinit/cgroups
native/etc/runlevels/sysinit/devfs
native/etc/runlevels/sysinit/dmesg
native/etc/runlevels/sysinit/sysfs
native/etc/runlevels/sysinit/udev
native/etc/runlevels/sysinit/udev-settle
native/etc/runlevels/sysinit/udev-trigger
native/etc/sysctl.d/README
native/lib/libeinfo.so
native/lib/libeinfo.so.1
native/lib/librc.so
native/lib/librc.so.1
native/lib/pkgconfig/einfo.pc
native/lib/pkgconfig/openrc.pc
native/sbin/halt
native/sbin/init
native/sbin/mountinfo
native/sbin/openrc
native/sbin/openrc-init
native/sbin/openrc-run
native/sbin/openrc-shutdown
native/sbin/poweroff
native/sbin/rc-service
native/sbin/rc-sstat
native/sbin/rc-update
native/sbin/reboot
native/sbin/shutdown
native/sbin/start-stop-daemon
native/sbin/supervise-daemon
native/usr/include/einfo.h
native/usr/include/rc.h
native/usr/libexec/rc/bin/checkpath
native/usr/libexec/rc/bin/ebegin
native/usr/libexec/rc/bin/eend
native/usr/libexec/rc/bin/eerror
native/usr/libexec/rc/bin/eerrorn
native/usr/libexec/rc/bin/eindent
native/usr/libexec/rc/bin/einfo
native/usr/libexec/rc/bin/einfon
native/usr/libexec/rc/bin/eoutdent
native/usr/libexec/rc/bin/esyslog
native/usr/libexec/rc/bin/eval_ecolors
native/usr/libexec/rc/bin/ewaitfile
native/usr/libexec/rc/bin/ewarn
native/usr/libexec/rc/bin/ewarnn
native/usr/libexec/rc/bin/ewend
native/usr/libexec/rc/bin/fstabinfo
native/usr/libexec/rc/bin/get_options
native/usr/libexec/rc/bin/is_newer_than
native/usr/libexec/rc/bin/is_older_than
native/usr/libexec/rc/bin/kill_all
native/usr/libexec/rc/bin/mountinfo
native/usr/libexec/rc/bin/on_ac_power
native/usr/libexec/rc/bin/rc-depend
native/usr/libexec/rc/bin/save_options
native/usr/libexec/rc/bin/service_crashed
native/usr/libexec/rc/bin/service_get_value
native/usr/libexec/rc/bin/service_hotplugged
native/usr/libexec/rc/bin/service_inactive
native/usr/libexec/rc/bin/service_set_value
native/usr/libexec/rc/bin/service_started
native/usr/libexec/rc/bin/service_started_daemon
native/usr/libexec/rc/bin/service_starting
native/usr/libexec/rc/bin/service_stopped
native/usr/libexec/rc/bin/service_stopping
native/usr/libexec/rc/bin/service_wasinactive
native/usr/libexec/rc/bin/shell_var
native/usr/libexec/rc/bin/vebegin
native/usr/libexec/rc/bin/veend
native/usr/libexec/rc/bin/veindent
native/usr/libexec/rc/bin/veinfo
native/usr/libexec/rc/bin/veoutdent
native/usr/libexec/rc/bin/vewarn
native/usr/libexec/rc/bin/vewend
native/usr/libexec/rc/sbin/mark_service_crashed
native/usr/libexec/rc/sbin/mark_service_failed
native/usr/libexec/rc/sbin/mark_service_hotplugged
native/usr/libexec/rc/sbin/mark_service_inactive
native/usr/libexec/rc/sbin/mark_service_started
native/usr/libexec/rc/sbin/mark_service_starting
native/usr/libexec/rc/sbin/mark_service_stopped
native/usr/libexec/rc/sbin/mark_service_stopping
native/usr/libexec/rc/sbin/mark_service_wasinactive
native/usr/libexec/rc/sbin/rc-abort
native/usr/libexec/rc/sbin/seedrng
native/usr/libexec/rc/sbin/swclock
native/usr/libexec/rc/sh/binfmt.sh
native/usr/libexec/rc/sh/cgroup-release-agent.sh
native/usr/libexec/rc/sh/functions.sh
native/usr/libexec/rc/sh/gendepends.sh
native/usr/libexec/rc/sh/init-early.sh
native/usr/libexec/rc/sh/init.sh
native/usr/libexec/rc/sh/openrc-run.sh
native/usr/libexec/rc/sh/rc-cgroup.sh
native/usr/libexec/rc/sh/rc-functions.sh
native/usr/libexec/rc/sh/rc-mount.sh
native/usr/libexec/rc/sh/runit.sh
native/usr/libexec/rc/sh/s6.sh
native/usr/libexec/rc/sh/start-stop-daemon.sh
native/usr/libexec/rc/sh/supervise-daemon.sh
native/usr/libexec/rc/version
native/usr/share/bash-completion/completions/openrc
native/usr/share/bash-completion/completions/openrc-service-script
native/usr/share/bash-completion/completions/rc-service
native/usr/share/bash-completion/completions/rc-status
native/usr/share/bash-completion/completions/rc-update
native/usr/share/man/man3/ebegin.3
native/usr/share/man/man3/ebeginn.3
native/usr/share/man/man3/ebeginv.3
native/usr/share/man/man3/ebeginvn.3
native/usr/share/man/man3/ebracket.3
native/usr/share/man/man3/eend.3
native/usr/share/man/man3/eendv.3
native/usr/share/man/man3/eerror.3
native/usr/share/man/man3/eerrorn.3
native/usr/share/man/man3/eerrorx.3
native/usr/share/man/man3/eindent.3
native/usr/share/man/man3/eindentv.3
native/usr/share/man/man3/einfo.3
native/usr/share/man/man3/einfon.3
native/usr/share/man/man3/einfov.3
native/usr/share/man/man3/einfovn.3
native/usr/share/man/man3/eoutdent.3
native/usr/share/man/man3/eoutdentv.3
native/usr/share/man/man3/eprefix.3
native/usr/share/man/man3/ewarn.3
native/usr/share/man/man3/ewarnn.3
native/usr/share/man/man3/ewarnv.3
native/usr/share/man/man3/ewarnvn.3
native/usr/share/man/man3/ewarnx.3
native/usr/share/man/man3/ewend.3
native/usr/share/man/man3/ewendv.3
native/usr/share/man/man3/rc_config.3
native/usr/share/man/man3/rc_config_list.3
native/usr/share/man/man3/rc_config_load.3
native/usr/share/man/man3/rc_config_value.3
native/usr/share/man/man3/rc_deptree.3
native/usr/share/man/man3/rc_deptree_depend.3
native/usr/share/man/man3/rc_deptree_depends.3
native/usr/share/man/man3/rc_deptree_free.3
native/usr/share/man/man3/rc_deptree_load.3
native/usr/share/man/man3/rc_deptree_order.3
native/usr/share/man/man3/rc_deptree_update.3
native/usr/share/man/man3/rc_deptree_update_needed.3
native/usr/share/man/man3/rc_find_pids.3
native/usr/share/man/man3/rc_plugin_hook.3
native/usr/share/man/man3/rc_runlevel.3
native/usr/share/man/man3/rc_runlevel_exists.3
native/usr/share/man/man3/rc_runlevel_get.3
native/usr/share/man/man3/rc_runlevel_list.3
native/usr/share/man/man3/rc_runlevel_set.3
native/usr/share/man/man3/rc_runlevel_starting.3
native/usr/share/man/man3/rc_runlevel_stopping.3
native/usr/share/man/man3/rc_service.3
native/usr/share/man/man3/rc_service_add.3
native/usr/share/man/man3/rc_service_daemon_set.3
native/usr/share/man/man3/rc_service_daemons_crashed.3
native/usr/share/man/man3/rc_service_delete.3
native/usr/share/man/man3/rc_service_description.3
native/usr/share/man/man3/rc_service_exists.3
native/usr/share/man/man3/rc_service_extra_commands.3
native/usr/share/man/man3/rc_service_in_runlevel.3
native/usr/share/man/man3/rc_service_mark.3
native/usr/share/man/man3/rc_service_plugable.3
native/usr/share/man/man3/rc_service_resolve.3
native/usr/share/man/man3/rc_service_schedule_clear.3
native/usr/share/man/man3/rc_service_schedule_start.3
native/usr/share/man/man3/rc_service_started_daemon.3
native/usr/share/man/man3/rc_service_state.3
native/usr/share/man/man3/rc_service_value_get.3
native/usr/share/man/man3/rc_service_value_set.3
native/usr/share/man/man3/rc_services_in_runlevel.3
native/usr/share/man/man3/rc_services_in_state.3
native/usr/share/man/man3/rc_services_scheduled.3
native/usr/share/man/man3/rc_services_scheduled_by.3
native/usr/share/man/man3/rc_stringlist.3
native/usr/share/man/man3/rc_stringlist_add.3
native/usr/share/man/man3/rc_stringlist_addu.3
native/usr/share/man/man3/rc_stringlist_delete.3
native/usr/share/man/man3/rc_stringlist_free.3
native/usr/share/man/man3/rc_stringlist_new.3
native/usr/share/man/man3/rc_stringlist_sort.3
native/usr/share/man/man3/rc_yesno.3
native/usr/share/man/man8/openrc-init.8
native/usr/share/man/man8/openrc-run.8
native/usr/share/man/man8/openrc-shutdown.8
native/usr/share/man/man8/openrc.8
native/usr/share/man/man8/rc-service.8
native/usr/share/man/man8/rc-sstat.8
native/usr/share/man/man8/rc-status.8
native/usr/share/man/man8/rc-update.8
native/usr/share/man/man8/start-stop-daemon.8
native/usr/share/man/man8/supervise-daemon.8
native/usr/share/openrc/support/deptree2dot/README.md
native/usr/share/openrc/support/deptree2dot/deptree2dot
native/usr/share/openrc/support/init.d.examples/README.md
native/usr/share/openrc/support/init.d.examples/avahi-dnsconfd
native/usr/share/openrc/support/init.d.examples/avahid
native/usr/share/openrc/support/init.d.examples/dbus
native/usr/share/openrc/support/init.d.examples/dhcpcd
native/usr/share/openrc/support/init.d.examples/hald
native/usr/share/openrc/support/init.d.examples/named
native/usr/share/openrc/support/init.d.examples/ntpd
native/usr/share/openrc/support/init.d.examples/openvpn
native/usr/share/openrc/support/init.d.examples/polkitd
native/usr/share/openrc/support/init.d.examples/sshd
native/usr/share/openrc/support/init.d.examples/wpa_supplicant
native/usr/share/openrc/support/openvpn/README.md
native/usr/share/openrc/support/openvpn/down.sh
native/usr/share/openrc/support/openvpn/up.sh
native/usr/share/openrc/support/sysvinit/README.md
native/usr/share/openrc/support/sysvinit/halt.sh
native/usr/share/openrc/support/sysvinit/inittab
native/usr/share/zsh/site-functions/_openrc
native/usr/share/zsh/site-functions/_rc-service
native/usr/share/zsh/site-functions/_rc-status
native/usr/share/zsh/site-functions/_rc-update
