# pkglinux-pkgsrc.2022Q2/net/samba4/files/init.d/nmbd
#!/sbin/openrc-run

command=/pkg/sbin/nmbd
pidfile=/run/nmbd.pid
command_args="-D"
name="NetBIOS name server"
extra_commands="reload"

depend() {
	provide samba
	use net
	before local
}

start_pre() {
	if ! [ -f /proc/cpuinfo ]; then
		echo "WARNING: Samba requires a Linux-compatible procfs!"
		echo "WARNING: Please mount /proc before starting Samba."
	fi
}
