# pkglinux-pkgsrc.2022Q2/native/boot/PLIST
@comment $NetBSD$
native/etc/init.d
native/etc/rc.d/init.d/checkfs
native/etc/rc.d/init.d/cleanfs
native/etc/rc.d/init.d/console
native/etc/rc.d/init.d/halt
native/etc/rc.d/init.d/localnet
native/etc/rc.d/init.d/modules
native/etc/rc.d/init.d/mountfs
native/etc/rc.d/init.d/mountvirtfs
native/etc/rc.d/init.d/network
native/etc/rc.d/init.d/rc
native/etc/rc.d/init.d/reboot
native/etc/rc.d/init.d/sendsignals
native/etc/rc.d/init.d/setclock
native/etc/rc.d/init.d/swap
native/etc/rc.d/init.d/sysctl
native/etc/rc.d/init.d/sysklogd
native/etc/rc.d/init.d/template
native/etc/rc.d/init.d/udev
native/etc/rc.d/init.d/udev_retry
native/etc/rc.d/rc0.d/K80network
native/etc/rc.d/rc0.d/K90sysklogd
native/etc/rc.d/rc0.d/S60sendsignals
native/etc/rc.d/rc0.d/S65swap
native/etc/rc.d/rc0.d/S70mountfs
native/etc/rc.d/rc0.d/S90localnet
native/etc/rc.d/rc0.d/S99halt
native/etc/rc.d/rc1.d/K80network
native/etc/rc.d/rc1.d/K90sysklogd
native/etc/rc.d/rc2.d/K80network
native/etc/rc.d/rc2.d/K90sysklogd
native/etc/rc.d/rc3.d/S10sysklogd
native/etc/rc.d/rc3.d/S20network
native/etc/rc.d/rc4.d/S10sysklogd
native/etc/rc.d/rc4.d/S20network
native/etc/rc.d/rc5.d/S10sysklogd
native/etc/rc.d/rc5.d/S20network
native/etc/rc.d/rc6.d/K80network
native/etc/rc.d/rc6.d/K90sysklogd
native/etc/rc.d/rc6.d/S60sendsignals
native/etc/rc.d/rc6.d/S65swap
native/etc/rc.d/rc6.d/S70mountfs
native/etc/rc.d/rc6.d/S90localnet
native/etc/rc.d/rc6.d/S99reboot
native/etc/rc.d/rcS.d/S00mountvirtfs
native/etc/rc.d/rcS.d/S05modules
native/etc/rc.d/rcS.d/S08localnet
native/etc/rc.d/rcS.d/S10udev
native/etc/rc.d/rcS.d/S20swap
native/etc/rc.d/rcS.d/S30checkfs
native/etc/rc.d/rcS.d/S40mountfs
native/etc/rc.d/rcS.d/S45cleanfs
native/etc/rc.d/rcS.d/S50udev_retry
native/etc/rc.d/rcS.d/S70console
native/etc/rc.d/rcS.d/S90sysctl
native/etc/sysconfig/createfiles
native/etc/sysconfig/modules
native/etc/sysconfig/rc.site
native/etc/sysconfig/udev_retry
native/lib/lsb
native/lib/services/init-functions
native/lib/services/ipv4-static
native/lib/services/ipv4-static-route
native/sbin/ifdown
native/sbin/ifup
native/usr/share/man/man8/ifdown.8
native/usr/share/man/man8/ifup.8
