# pkglinux-pkgsrc.2022Q2/linux/flatpak/PLIST
@comment $NetBSD$
bin/flatpak
bin/flatpak-bisect
bin/flatpak-coredumpctl
etc/dbus-1/system.d/org.freedesktop.Flatpak.SystemHelper.conf
etc/profile.d/flatpak.sh
include/flatpak/flatpak-bundle-ref.h
include/flatpak/flatpak-enum-types.h
include/flatpak/flatpak-error.h
include/flatpak/flatpak-installation.h
include/flatpak/flatpak-installed-ref.h
include/flatpak/flatpak-instance.h
include/flatpak/flatpak-portal-error.h
include/flatpak/flatpak-ref.h
include/flatpak/flatpak-related-ref.h
include/flatpak/flatpak-remote-ref.h
include/flatpak/flatpak-remote.h
include/flatpak/flatpak-transaction.h
include/flatpak/flatpak-version-macros.h
include/flatpak/flatpak.h
lib/girepository-1.0/Flatpak-1.0.typelib
lib/libflatpak.la
lib/pkgconfig/flatpak.pc
lib/systemd/system-environment-generators/60-flatpak-system-only
lib/systemd/system/flatpak-system-helper.service
lib/systemd/user-environment-generators/60-flatpak
lib/systemd/user/flatpak-oci-authenticator.service
lib/systemd/user/flatpak-portal.service
lib/systemd/user/flatpak-session-helper.service
lib/sysusers.d/flatpak.conf
libexec/flatpak-bwrap
libexec/flatpak-dbus-proxy
libexec/flatpak-oci-authenticator
libexec/flatpak-portal
libexec/flatpak-session-helper
libexec/flatpak-system-helper
libexec/flatpak-validate-icon
libexec/revokefs-fuse
share/bash-completion/completions/flatpak
share/dbus-1/interfaces/org.freedesktop.Flatpak.Authenticator.xml
share/dbus-1/interfaces/org.freedesktop.Flatpak.xml
share/dbus-1/interfaces/org.freedesktop.portal.Flatpak.xml
share/dbus-1/services/org.flatpak.Authenticator.Oci.service
share/dbus-1/services/org.freedesktop.Flatpak.service
share/dbus-1/services/org.freedesktop.portal.Flatpak.service
share/dbus-1/system-services/org.freedesktop.Flatpak.SystemHelper.service
share/fish/vendor_completions.d/flatpak.fish
share/fish/vendor_conf.d/flatpak.fish
share/flatpak/triggers/desktop-database.trigger
share/flatpak/triggers/gtk-icon-cache.trigger
share/flatpak/triggers/mime-database.trigger
share/gir-1.0/Flatpak-1.0.gir
share/locale/cs/LC_MESSAGES/flatpak.mo
share/locale/da/LC_MESSAGES/flatpak.mo
share/locale/de/LC_MESSAGES/flatpak.mo
share/locale/en_GB/LC_MESSAGES/flatpak.mo
share/locale/es/LC_MESSAGES/flatpak.mo
share/locale/fr/LC_MESSAGES/flatpak.mo
share/locale/gl/LC_MESSAGES/flatpak.mo
share/locale/hi/LC_MESSAGES/flatpak.mo
share/locale/hr/LC_MESSAGES/flatpak.mo
share/locale/hu/LC_MESSAGES/flatpak.mo
share/locale/id/LC_MESSAGES/flatpak.mo
share/locale/oc/LC_MESSAGES/flatpak.mo
share/locale/pl/LC_MESSAGES/flatpak.mo
share/locale/pt/LC_MESSAGES/flatpak.mo
share/locale/pt_BR/LC_MESSAGES/flatpak.mo
share/locale/ro/LC_MESSAGES/flatpak.mo
share/locale/ru/LC_MESSAGES/flatpak.mo
share/locale/sk/LC_MESSAGES/flatpak.mo
share/locale/sv/LC_MESSAGES/flatpak.mo
share/locale/tr/LC_MESSAGES/flatpak.mo
share/locale/uk/LC_MESSAGES/flatpak.mo
share/locale/zh_CN/LC_MESSAGES/flatpak.mo
share/locale/zh_TW/LC_MESSAGES/flatpak.mo
share/polkit-1/actions/org.freedesktop.Flatpak.policy
share/polkit-1/rules.d/org.freedesktop.Flatpak.rules
share/zsh/site-functions/_flatpak
