# pkglinux-pkgsrc.2022Q2/linux/shellinabox/PLIST.0
@comment $NetBSD$
bin/shellinaboxd
man/man1/shellinaboxd.1
share/doc/shellinabox/AUTHORS
share/doc/shellinabox/COPYING
share/doc/shellinabox/GPL-2
share/doc/shellinabox/NEWS
share/shellinaboxd/black-on-white.css
share/shellinaboxd/white-on-black.css
share/examples/rc.d/shellinaboxd
