# pkglinux-pkgsrc.2022Q2/native/bootscripts/PLIST
@comment $NetBSD$
native/boot/grub/grub.cfg
native/boot/syslinux/syslinux.cfg
native/etc/init.d
native/etc/inittab.alex
native/etc/inittab.bsd
native/etc/inittab.lfs
native/etc/rc.d/rc0
native/etc/rc.d/rc1
native/etc/rc.d/rc2
native/etc/rc.d/rc3
native/etc/rc.d/rc4
native/etc/rc.d/rc5
native/etc/rc.d/rc6
native/etc/rc.d/rcS
native/etc/sysconfig/console
native/sbin/lndir
native/sbin/lndir.sh
native/sbin/pkg-chroot
native/sbin/pkglinux-install
native/sbin/pkglinux-makeall
native/sbin/pkglinux-makeiso
native/sbin/pkglinux-makeusb
native/sbin/pkglinux-upgrade
native/sbin/pkg-strap
native/sbin/pkg-sum
native/sbin/powersave
native/sbin/service
native/sbin/wpa.sh
native/usr/bin/which
native/usr/share/man/man1/pkg-chroot.1
native/usr/share/man/man1/pkglinux-install.1
native/usr/share/man/man1/pkglinux-makeall.1
native/usr/share/man/man1/pkglinux-makeiso.1
native/usr/share/man/man1/pkglinux-makeusb.1
native/usr/share/man/man1/pkglinux-upgrade.1
native/usr/share/man/man1/pkg-strap.1
@pkgdir native/etc/rc.d/rcS.d
@pkgdir native/etc/rc.d/rc6.d
@pkgdir native/etc/rc.d/rc5.d
@pkgdir native/etc/rc.d/rc4.d
@pkgdir native/etc/rc.d/rc3.d
@pkgdir native/etc/rc.d/rc2.d
@pkgdir native/etc/rc.d/rc1.d
@pkgdir native/etc/rc.d/rc0.d
@pkgdir native/etc/rc.d/init.d
